# dynaruby.rcd
#!/bin/ksh

daemon="/usr/local/sbin/dynaruby"
name="dynaruby"

. /etc/rc.d/rc.subr

rc_start() {
export DYNARUBY_KEY="YOUR,KEY"
exec /usr/local/sbin/dynaruby &
}
rc_cmd $1