# dynaruby.rcd
#!/bin/ksh

daemon="/usr/local/sbin/dynaruby"
name="dynaruby"

. /etc/rc.d/rc.subr

rc_start() {
export DYNARUBY_KEY="YOUR,KEY"
time=$(date '+%H:%M:%S')
echo "$time Starting $name" >> /var/log/dynaruby.log
exec $daemon &
}
rc_cmd $1