# dynaruby.rcd
#!/bin/ksh

daemon="/usr/local/sbin/dynaruby"
name="dynaruby"

. /etc/rc.d/rc.subr

export DYNARUBY_KEY=''

rc_cmd $1