# dynaruby.rcd
#!/bin/ksh

daemon="/usr/local/sbin/dynaruby"
name="dynaruby"

. /etc/rc.d/rc.subr

rc_start()
{
    var DYNARUBY_KEY="1234" ${daemon}
}

rc_cmd $1